# beagle/startup-service/media.automount
[Unit]
Description=Automount the SD Card

[Automount]
Where=/media
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target

